# src/DEBIAN/postinst
#! /bin/bash

if [ "$1" = "triggered" ]
then
	if [ "$2" = "/opt/ont/minknow/ont-python/lib" ]
	then
		logger "fast5upload: Triggered for ONT Python"
	fi
fi

mkdir -p /var/lib/mlstverse
chown root:root /var/lib/mlstverse

/opt/ont/minknow/ont-python/bin/python3 -m pip install watchdog

if [ $(head -c 1 /etc/mlstverse/fast5upload.conf) == '{' ]
then
	printf 'You are using an old version of fast5upload config.\n'
	printf 'We will attempt to fix the config format automatically the next time you start fast5upload.\n'
	printf 'If anything did not work as expected, please contact support.\n'
fi

# Auto-generated script copy from deb-systemd-helper for systemd unit management
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	systemctl --system daemon-reload
        # This will only remove masks created by d-s-h on package removal.
        deb-systemd-helper unmask 'fast5upload.service' >/dev/null || true

        # was-enabled defaults to true, so new installations run enable.
        if deb-systemd-helper --quiet was-enabled 'fast5upload.service'; then
                # Enables the unit on first installation, creates new
                # symlinks on upgrades if the unit file has changed.
                deb-systemd-helper enable 'fast5upload.service' >/dev/null || true
        else
                # Update the statefile to add new symlinks (if any), which need to be
                # cleaned up on purge. Also remove old symlinks.
                deb-systemd-helper update-state 'fast5upload.service' >/dev/null || true
        fi
fi

